# .config/systemd/user/alarm2.service
[Unit]
Description="Alarms"

[Service]
Type=oneshot
ExecStart="/home/moeru/alarms/play_alarm2.sh"

[Install]
WantedBy=default.target


