# systemd/user/alarm.service
[Unit]
Description="Alarms"

[Service]
Type=oneshot
ExecStart="/home/moeru/alarms/play_alarm.sh"

[Install]
WantedBy=default.target


