# systemd/user/helloworld.service
[Unit]
Description="Hello World script"

[Service]
Type=oneshot
ExecStart=/home/moeru/helloworld.sh

[Install]
WantedBy=default.target


