# systemd/user/helloworld.timer
[Unit]
Description=Run "Do some thing" every 1 minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
